# web_service_restart.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/tmp/webapp/packer/wb_restart.sh

[Install]
WantedBy=multi-user.target
